# copy-github-to-github.service
[Unit]
Description=copy-github-to-github

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=$CMD

[Install]
WantedBy=multi-user.target
